# systemd/supermon-ng-database-update.service
[Unit]
Description=Supermon-NG Database Auto-Update Service
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/html/supermon-ng
ExecStart=/usr/bin/php /var/www/html/supermon-ng/scripts/database-auto-update.php

[Install]
WantedBy=multi-user.target

