# systemd/supermon-ng-node-status.service
[Unit]
Description=Supermon-NG Node Status Update Service
After=network.target asterisk.service
Wants=asterisk.service

[Service]
Type=oneshot
User=root
WorkingDirectory=APP_DIR_PLACEHOLDER/user_files/sbin
ExecStart=/usr/bin/python3 APP_DIR_PLACEHOLDER/user_files/sbin/ast_node_status_update.py
StandardOutput=append:APP_DIR_PLACEHOLDER/logs/node-status-update.log
StandardError=append:APP_DIR_PLACEHOLDER/logs/node-status-update.log

[Install]
WantedBy=multi-user.target
