# debian/postinst
#!/bin/sh
POSTINST_MODE="${1-}"
PREVIOUS_VERSION="${2-}"

APP_DIR=/var/www/html/supermon-ng
WWW_USER=www-data
WWW_GROUP=www-data
NODE_STATUS_TIMER_DROPIN=/etc/systemd/system/supermon-ng-node-status.timer.d/interval.conf

load_env() {
    APP_BASE_PATH="${APP_BASE_PATH:-/supermon-ng}"
    NODE_STATUS_INTERVAL_MINUTES="${NODE_STATUS_INTERVAL_MINUTES:-5}"

    if [ -f "$APP_DIR/.env" ]; then
        # shellcheck disable=SC1091
        set -a
        . "$APP_DIR/.env" 2>/dev/null || true
        set +a
        APP_BASE_PATH="${APP_BASE_PATH:-/supermon-ng}"
        NODE_STATUS_INTERVAL_MINUTES="${NODE_STATUS_INTERVAL_MINUTES:-5}"
    fi
}

read_debconf_bool() {
    key="$1"
    default="$2"

    if ! command -v debconf-communicate >/dev/null 2>&1; then
        echo "$default"
        return 0
    fi

    ans=$(echo "GET $key" | debconf-communicate supermon-ng 2>/dev/null | head -1)
    case "$ans" in
        "0 true") echo true ;;
        "0 false") echo false ;;
        *) echo "$default" ;;
    esac
}

apply_app_base_path() {
    load_env
    if [ -f "$APP_DIR/.env" ] && [ -x "$APP_DIR/scripts/configure-app-base-path.sh" ]; then
        export APP_DIR APP_BASE_PATH
        "$APP_DIR/scripts/configure-app-base-path.sh" "$APP_DIR" || true
    fi
}

ensure_runtime_dirs() {
    mkdir -p "$APP_DIR/logs" "$APP_DIR/cache" "$APP_DIR/database"
    chown -R "${WWW_USER}:${WWW_GROUP}" "$APP_DIR/logs" "$APP_DIR/cache" "$APP_DIR/database" 2>/dev/null || true
    chmod 755 "$APP_DIR/logs" "$APP_DIR/cache" "$APP_DIR/database" 2>/dev/null || true
}

ensure_app_permissions() {
    chown -R "${WWW_USER}:${WWW_GROUP}" "$APP_DIR" 2>/dev/null || true
    chmod -R 755 "$APP_DIR" 2>/dev/null || true
    chmod -R 755 "$APP_DIR/logs" "$APP_DIR/cache" 2>/dev/null || true

    find "$APP_DIR/user_files" -type d -exec chmod 755 {} \; 2>/dev/null || true
    find "$APP_DIR/user_files" -type f -exec chmod 644 {} \; 2>/dev/null || true
    chmod 644 "$APP_DIR/user_files/.htaccess" 2>/dev/null || true
    chmod 644 "$APP_DIR/user_files/.htpasswd" 2>/dev/null || true

    for script in ast_node_status_update.py din ssinfo dvswitch-bridge-restart.sh; do
        [ -f "$APP_DIR/user_files/sbin/$script" ] && chmod 755 "$APP_DIR/user_files/sbin/$script"
    done
    [ -f "$APP_DIR/user_files/sbin/node_info.ini" ] && chmod 644 "$APP_DIR/user_files/sbin/node_info.ini"

    if [ -f "$APP_DIR/.env" ]; then
        chown root:${WWW_GROUP} "$APP_DIR/.env" 2>/dev/null || true
        chmod 640 "$APP_DIR/.env" 2>/dev/null || true
    fi
}

validate_sudoers() {
    if [ -f /etc/sudoers.d/011-supermon-ng ]; then
        visudo -cf /etc/sudoers.d/011-supermon-ng >/dev/null 2>&1 || {
            echo "supermon-ng: sudoers validation failed for /etc/sudoers.d/011-supermon-ng" >&2
            exit 1
        }
    fi
}

apply_node_status_timer_interval() {
    load_env
    interval="${NODE_STATUS_INTERVAL_MINUTES:-5}"

    mkdir -p "$(dirname "$NODE_STATUS_TIMER_DROPIN")"
    cat >"$NODE_STATUS_TIMER_DROPIN" <<EOF
# Managed by supermon-ng postinst from NODE_STATUS_INTERVAL_MINUTES in .env
[Unit]
Description=Run Supermon-NG Node Status Update every ${interval} minutes

[Timer]
OnUnitActiveSec=${interval}min
EOF
    chmod 644 "$NODE_STATUS_TIMER_DROPIN"
}

generate_local_allmon_if_fresh() {
    if [ -n "$PREVIOUS_VERSION" ]; then
        return 0
    fi

    if [ ! -f "$APP_DIR/scripts/generate_local_allmon.php" ]; then
        return 0
    fi

    echo "supermon-ng: generating user_files/allmon.ini from Asterisk config..."
    php "$APP_DIR/scripts/generate_local_allmon.php" --force || {
        echo "supermon-ng: allmon.ini generation failed; edit user_files/allmon.ini manually or run:" >&2
        echo "  sudo php $APP_DIR/scripts/generate_local_allmon.php --force" >&2
    }
}

configure_apache_from_debconf() {
    configure_apache=$(read_debconf_bool supermon-ng/configure-apache true)
    disable_default_sites=$(read_debconf_bool supermon-ng/disable-default-sites true)
    refresh_apache_site=$(read_debconf_bool supermon-ng/refresh-apache-site false)

    SKIP_APACHE=false
    DISABLE_DEFAULT_SITES=true
    OVERWRITE_SITE=false

    if [ "$configure_apache" = "false" ]; then
        SKIP_APACHE=true
    fi
    if [ "$disable_default_sites" = "false" ]; then
        DISABLE_DEFAULT_SITES=false
    fi
    if [ "$refresh_apache_site" = "true" ]; then
        OVERWRITE_SITE=true
    fi

    if [ ! -x "$APP_DIR/scripts/configure-apache.sh" ]; then
        echo "supermon-ng: configure-apache.sh missing; Apache setup skipped" >&2
        return 0
    fi

    export APP_DIR SKIP_APACHE DISABLE_DEFAULT_SITES OVERWRITE_SITE
    export ALLOW_APT_INSTALL=false CONFIGURE_LOG_ACLS=true
    "$APP_DIR/scripts/configure-apache.sh" configure || true
}

enable_systemd_units() {
    systemctl daemon-reload || true
    systemctl enable supermon-ng-backend 2>/dev/null || true
    systemctl enable supermon-ng-websocket.service 2>/dev/null || true
    systemctl enable supermon-ng-node-status.timer 2>/dev/null || true
    systemctl enable supermon-ng-database-update.timer 2>/dev/null || true
    systemctl start supermon-ng-backend 2>/dev/null || true
    systemctl start supermon-ng-websocket.service 2>/dev/null || true
    systemctl start supermon-ng-node-status.timer 2>/dev/null || true
    systemctl start supermon-ng-database-update.timer 2>/dev/null || true
}

case "$POSTINST_MODE" in
    configure)
        if [ ! -f "$APP_DIR/.env" ] && [ -f "$APP_DIR/.env.example" ]; then
            cp "$APP_DIR/.env.example" "$APP_DIR/.env"
            chown root:${WWW_GROUP} "$APP_DIR/.env" 2>/dev/null || true
            chmod 640 "$APP_DIR/.env" 2>/dev/null || true
        fi

        ensure_runtime_dirs
        apply_app_base_path
        generate_local_allmon_if_fresh
        validate_sudoers
        apply_node_status_timer_interval
        configure_apache_from_debconf
        enable_systemd_units

        if [ ! -d "$APP_DIR/vendor" ] && [ -x "$APP_DIR/scripts/composer-install-production.sh" ]; then
            if command -v composer >/dev/null 2>&1; then
                "$APP_DIR/scripts/composer-install-production.sh" "$APP_DIR" || true
            fi
        fi

        ensure_app_permissions
        ;;

    abort-upgrade|abort-remove|abort-deconfigure|triggered)
        ;;

    "")
        echo "supermon-ng: postinst called with empty argument; skipping custom setup" >&2
        ;;

    *)
        echo "supermon-ng: postinst skipping custom setup for argument \`$POSTINST_MODE'" >&2
        ;;
esac

#DEBHELPER#

exit 0
